# config/dhcpv6-ipset.service
[Unit]
Description=DHCPv6 IPSet Service

[Service]
ExecStart=/config/scripts/dhcpv6-ipset

[Install]
WantedBy=multi-user.target
